# dispatch_listener/rootfs/etc/services.d/dispatch_listener/run
#!/usr/bin/env sh
set -e
cd /opt/dispatch_listener
exec python3 -m main
